# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx_err.log

[program:phase1]
command=python3 -m uvicorn phase1_knowledgeBase.app.main:app --host 0.0.0.0 --port 8101
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/phase1.log
stderr_logfile=/var/log/phase1_err.log

[program:phase2]
command=python3 -m uvicorn phase2_weeklyPulse.app.main:app --host 0.0.0.0 --port 8102
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/phase2.log
stderr_logfile=/var/log/phase2_err.log

[program:phase3]
command=python3 -m uvicorn phase3_voiceScheduler.app.main:app --host 0.0.0.0 --port 8103
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/phase3.log
stderr_logfile=/var/log/phase3_err.log

[program:phase4]
command=python3 -m uvicorn phase4_integrationHub.app.main:app --host 0.0.0.0 --port 8000
directory=/app
environment=PHASE1_BASE_URL="http://localhost:8101",PHASE2_BASE_URL="http://localhost:8102",PHASE3_BASE_URL="http://localhost:8103"
autostart=true
autorestart=true
startsecs=8
stdout_logfile=/var/log/phase4.log
stderr_logfile=/var/log/phase4_err.log
