# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx_err.log

[program:phase1]
command=python3 -m uvicorn phase1_knowledgeBase.app.main:app --host 0.0.0.0 --port 8101
directory=/app
autostart=true
autorestart=true
environment=PINECONE_API_KEY="%(ENV_PINECONE_API_KEY)s",PINECONE_INDEX_NAME="%(ENV_PINECONE_INDEX_NAME)s",PINECONE_CLOUD="%(ENV_PINECONE_CLOUD)s",PINECONE_REGION="%(ENV_PINECONE_REGION)s",OPENROUTER_API_KEY="%(ENV_OPENROUTER_API_KEY)s",OPENROUTER_BASE_URL="%(ENV_OPENROUTER_BASE_URL)s",OPENROUTER_CHAT_MODEL="%(ENV_OPENROUTER_CHAT_MODEL)s"
stdout_logfile=/var/log/phase1.log
stderr_logfile=/var/log/phase1_err.log

[program:phase2]
command=python3 -m uvicorn phase2_weeklyPulse.app.main:app --host 0.0.0.0 --port 8102
directory=/app
autostart=true
autorestart=true
environment=GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",GROQ_BASE_URL="%(ENV_GROQ_BASE_URL)s"
stdout_logfile=/var/log/phase2.log
stderr_logfile=/var/log/phase2_err.log

[program:phase3]
command=python3 -m uvicorn phase3_voiceScheduler.app.main:app --host 0.0.0.0 --port 8103
directory=/app
autostart=true
autorestart=true
environment=SARVAM_API_KEY="%(ENV_SARVAM_API_KEY)s",GOOGLE_TOKEN_JSON="%(ENV_GOOGLE_TOKEN_JSON)s",GOOGLE_TRACKING_DOC_ID="%(ENV_GOOGLE_TRACKING_DOC_ID)s",ADVISOR_EMAIL="%(ENV_ADVISOR_EMAIL)s",GOOGLE_MEET_URL="%(ENV_GOOGLE_MEET_URL)s",MCP_SERVER_URL="%(ENV_MCP_SERVER_URL)s",PHASE2_BASE_URL="http://localhost:8102"
stdout_logfile=/var/log/phase3.log
stderr_logfile=/var/log/phase3_err.log

[program:phase4]
command=python3 -m uvicorn phase4_integrationHub.app.main:app --host 0.0.0.0 --port 8000
directory=/app
environment=PHASE1_BASE_URL="http://localhost:8101",PHASE2_BASE_URL="http://localhost:8102",PHASE3_BASE_URL="http://localhost:8103"
autostart=true
autorestart=true
startsecs=8
stdout_logfile=/var/log/phase4.log
stderr_logfile=/var/log/phase4_err.log
